# org.glite.lb.server/config/glite-lb-proxy-purge.cron
31 1 * * *	glite	( export GLITE_LB_EXPORT_BKSERVER=`hostname -f`; export GLITE_LB_PURGE_TARGET_RUNTIME=84600s; export GLITE_LB_PURGE_OTHER_OPTIONS="$GLITE_LB_PURGE_OTHER_OPTIONS -x"; printf "`date "+\%b \%e \%T"` `hostname -s` glite-lb-export: "; @glite_prefix@/sbin/glite-lb-export ) >> /var/log/glite/glite-lb-server-purge-proxy.log 2>&1
