# org.glite.lb.server/config/glite-lb-notif-keeper.cron
0 */4 * * *	glite	. /etc/profile.d/grid-env.sh ; export GLITE_WMS_NOTIF_SERVER=${HOSTNAME}:9000; @glite_prefix@/sbin/glite-lb-notif-keeper --site-notif /etc/glite-lb/site-notif.conf > /tmp/notif.keeper.log.$$; logger -t glite-lb-notif-keeper -f /tmp/notif.keeper.log.$$ 
