# org.glite.lb.server/config/glite-lb-notif-keeper.cron
0 */4 * * *	glite	export GLITE_WMS_NOTIF_SERVER=`hostname -f`:9000; @glite_prefix@/sbin/glite-lb-notif-keeper --file-prefix /var/spool/glite/lb-notif/dglogd.log --site-notif @glite_etc@/glite-lb/site-notif.conf > /tmp/notif.keeper.log.$$; logger -t glite-lb-notif-keeper -f /tmp/notif.keeper.log.$$; rm -f /tmp/notif.keeper.log.$$
